# docker/cron/crontab
* * * * * cd /app && sudo -u www-data /usr/local/bin/php artisan schedule:run >> /var/log/cron.log
